# uwsgi-portal/install/etc/services.d/uwsgi/run
#!/command/with-contenv sh

set -e

/usr/bin/uwsgi --yaml /etc/uwsgi/apps-available/llng-server.yaml ${UWSGI_ARGS}
